# awstats.crontab
05 * * * * stats umask 002; /usr/bin/awstats_updateall.pl now -configdir=/etc/webapps/awstats -awstatsprog=/usr/bin/awstats.pl > /dev/null
# consider also running what is above from /etc/logrotate.d/httpd as prerotate script - see:
# http://awstats.sourceforge.net/docs/awstats_faq.html#ROTATE
